# mma-mcp.service
[Unit]
Description=mma-mcp — Wolfram Engine MCP Server
Documentation=https://github.com/your-org/mma-mcp
After=network.target

[Service]
Type=simple
User=mma
Group=mma
WorkingDirectory=/opt/mma-mcp

# Use the project's virtual environment
ExecStart=/opt/mma-mcp/.venv/bin/mma-mcp serve --transport http --host 127.0.0.1 --port 8000

# Environment variables for auth token or DNS provider credentials.
# Prefer EnvironmentFile for secrets (permissions 0600, owned by root).
# EnvironmentFile=/opt/mma-mcp/.env

Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# Allow write to the working directory (for logs, temp files)
ReadWritePaths=/opt/mma-mcp
# Wolfram Engine needs access to its installation and license files
ReadOnlyPaths=/usr/local/Wolfram
ReadWritePaths=/home/mma/.Mathematica
ReadWritePaths=/home/mma/.WolframEngine

[Install]
WantedBy=multi-user.target
